# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/mysqld_safe

[program:springapp]
command=java -jar /workspace/app/build/libs/spring-docker.jar
